# src/crontab.d/send_notifications.cron
*/5 * * * * php <>/bin/console app:send-notification
