# src/crontab.d/send_notifications.cron
*/5 * * * * docker exec -it sggwmeetapp_dockerized_php_1 php /var/www/symfony/bin/console app:send-notification
