# src/crontab.d/send_notifications.cron
*/5 * * * * docker exec -it sggwmeetapp_dockerized-php-1 php /var/www/symfony/bin/console app:send-notification
